# bazzite/bazzite-hud-svg/bazzite-hud.user.timer
# ==============================================================================
# Bazzite HUD SVG Generator
# Generates an SVG displaying rpm-ostree status for desktop widgets.
# Source: https://github.com/Jahfry/Immutable-SteamOS-Bazzite-etc/tree/main/bazzite/bazzite-hud-svg
# ==============================================================================

[Unit]
Description=Bazzite HUD Timer (User)

[Timer]
# Run on boot, then every hour
OnStartupSec=10sec
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
