# bazzite/bazzite-hud-svg/bazzite-hud-trigger.system.service
# ==============================================================================
# Bazzite HUD SVG Generator
# Generates an SVG displaying rpm-ostree status for desktop widgets.
# Source: https://github.com/Jahfry/Immutable-SteamOS-Bazzite-etc/tree/main/bazzite/bazzite-hud-svg
# ==============================================================================

[Unit]
Description=Bazzite HUD Generator (System)

[Service]
Type=oneshot
User=bazzite
# Critical: Explicitly set HOME so the script finds the correct path during boot
Environment="HOME=/var/home/bazzite"
ExecStart=/usr/bin/bash /var/home/bazzite/.local/bin/bazzite-hud-svg.sh

